# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/run/supervisord/supervisord.pid

[unix_http_server]
file = /run/supervisord/supervisor.sock
chmod = 0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:prometheus]
directory=/var/lib/prometheus
command=/usr/bin/prometheus --config.file=/etc/prometheus/prometheus.yml
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:grafana]
directory=/usr/share/grafana
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana --config=/etc/grafana.ini --packaging=docker cfg:default.log.mode="console"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
